# etc/systemd/system/pwrmgr-restarter.path
[Path]
PathChanged=/home/pi/.openauto/scripts/pwrmgr.json

[Install]
WantedBy=multi-user.target